# gguf_server.service
[Unit]
Description=GGUF Audio Server (llama-liquid-audio-server)
After=network.target

[Service]
Type=simple
User=audioservice
WorkingDirectory=/home/audioservice
ExecStart=/home/audioservice/gguf-runners/llama-liquid-audio-ubuntu-x64/llama-liquid-audio-server -m /home/audioservice/models/LFM2.5-Audio-1.5B-GGUF/LFM2.5-Audio-1.5B-Q4_0.gguf -mm /home/audioservice/models/LFM2.5-Audio-1.5B-GGUF/mmproj-LFM2.5-Audio-1.5B-Q4_0.gguf -mv /home/audioservice/models/LFM2.5-Audio-1.5B-GGUF/vocoder-LFM2.5-Audio-1.5B-Q4_0.gguf --tts-speaker-file /home/audioservice/models/LFM2.5-Audio-1.5B-GGUF/tokenizer-LFM2.5-Audio-1.5B-Q4_0.gguf --host 127.0.0.1 --port 8080
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
